# bridge/companion-bridge.service
# systemd unit for the Cooking Companion subscription bridge.
# Install: sudo cp companion-bridge.service /etc/systemd/system/
#          sudo systemctl enable --now companion-bridge
[Unit]
Description=Cook Anything companion bridge (headless Claude Code)
After=network-online.target
Wants=network-online.target
# Phase 0 safety interlock. bridge/phase0-disable.sh creates this marker.
ConditionPathExists=!/opt/cook-anything/BRIDGE_DISABLED

[Service]
Type=simple
User=companion
WorkingDirectory=/opt/cook-anything/bridge
# BRIDGE_TOKEN, CLAUDE_CODE_OAUTH_TOKEN, COMPANION_MODEL, CLAUDE_BIN live here:
EnvironmentFile=/etc/companion-bridge.env
ExecStart=/usr/bin/env node server.mjs
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target
